# src/bing/crontab
00 */1 * * * /usr/bin/python /path/to/github.com/typhoonzero/dailywall4mac/src/bing/set_desktop_picture_of_the_day.py 
